# systemd/quadlet.container
[Unit]
Description=

[Container]
Image=

[Service]

[Install]
WantedBy=
